# service-hub_1.0.0_arm/DEBIAN/postinst
#!/bin/bash
set -e

systemctl daemon-reload
systemctl enable service_hubd.service
systemctl restart service_hubd.service

exit 0

